# crontab
# min hour dom month dow command

0 */6 * * *  /home/siim/bithumb/scrape_and_push.sh >> /tmp/scrape.log
